# .github/docker/db/etc-supervisor-conf.d-services.conf
[program:mysqld]
#command = /usr/bin/mysqld_safe
command = /etc/init.d/mysql start
user = root
autostart = true