# .github/docker/ptserver/etc-supervisor-conf.d-services.conf
[program:node]
directory = /gt/sc-prog-repos/ptfile/ptserver
command = /bin/bash -c "npm install && npm run start-test"
user = root
autostart = true