# .github/docker/puppet/etc-supervisor-conf.d-services.conf
[program:node]
directory = /gt/sc-prog-repos/ptfile/puppet
command = /bin/bash 
user = root
autostart = true
autorestart=false


